# recipes-connectivity/lte-mbim-conf/files/init_quectel_mbim.service
[Unit]
Description=Init quectel mbim modem

[Service]
Type=simple
ExecStart=/usr/bin/init.quectel.mbim.sh start
ExecStop=/usr/bin/init.quectel.mbim.sh stop
Restart=on-failure

[Install]
WantedBy=multi-user.target
